# supervisord/script/conf.d/airflow.ini
[program:airflow_webserver]
command=airflow webserver
user=airflow
redirect_stderr=true
stdout_logfile=/data/applogs/airflow/web.log
autostart=true
autorestart=false
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
environment=AIRFLOW_HOME=/data/appdata/airflow
directory=/data/appdata/airflow

[program:airflow_scheduler]
command=airflow scheduler
user=airflow
redirect_stderr=true
stdout_logfile=/data/applogs/airflow/scheduler.log
autostart=true
autorestart=false
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
environment=AIRFLOW_HOME=/data/appdata/airflow
directory=/data/appdata/airflow

[program:airflow_worker]
command=airflow worker
user=airflow
redirect_stderr=true
stdout_logfile=/data/applogs/airflow/worker.log
autostart=true
autorestart=false
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
environment=AIRFLOW_HOME=/data/appdata/airflow
directory=/data/appdata/airflow
